# ansible/roles/watchtimer/templates/watchlog.timer
[Unit]
Description=Timer for Watchlog

[Timer]
OnUnitActiveSec=30

[Install]
WantedBy=multi-user.target
